# conf_cnz/deploy/etc/systemd/system/dhcpd.service.d/lde.conf
[Unit]
ConditionPathExists=/etc/cluster/nodes/this/properties/bootserver

[Service]

