# conf_cnz/deploy/usr/lib/systemd/system/apg-ldap.service
[Unit]
Description=APG LDAP

[Service]
ExecStart=/usr/lib/systemd/scripts/apg-ldap.sh
Type=forking
Restart=always
