# conf_cnz/deploy/usr/lib/systemd/system/apg-rsh@.service
[Unit]
Description=Remote Shell Facilities Server
After=lde-servicestart.target

[Service]
ExecStart=-/usr/sbin/in.rshd
StandardInput=socket
