# conf_cnz/deploy/usr/lib/systemd/system/apg-vsftpd.socket
[Unit]
Conflicts=vsftpd.service
DefaultDependencies=no

[Socket]
ListenStream=<node_public_network_ip_address>:21
Accept=yes

[Install]
WantedBy=lde-servicestart.target
