# conf_cnz/deploy/usr/lib/systemd/system/apg-netconf-beep@.service
[Unit]
Description=NETCONF-beep Server
After=local-fs.target

[Service]
User=netcc
ExecStart=-/opt/com/bin/netconf
StandardInput=socket
