# conf_cnz/deploy/usr/lib/systemd/system/dhcpd.service
[Unit]
Description=ISC DHCPv4 Server

[Service]
ExecStart=/usr/bin/true
Type=oneshot
RemainAfterExit=no
