# conf_cnz/deploy/files
cluster/hooks/after-booting-from-disk.tar.gz|444
cluster/hooks/post-installation.tar.gz|444
cluster/hooks/pre-installation.tar.gz|444
etc/audit/rules.d/901-apg-users.rules|444
etc/bindresvport.blacklist|644
etc/dhcpd.conf.local|755
etc/dhcpd.conf.local_vm|755
etc/10-apg-rsyslog-rule|644
etc/02-lde-syslog-logstream-list.conf|644
etc/ntp.conf.local|755
etc/chrony.conf.local|755
etc/init.d/lde-network|755
etc/nscd.conf|644
etc/nsswitch.conf|644
etc/ssh_key_rp|640
etc/openldap/schema/euac-extended.schema|755
etc/openldap/slapd.conf|640
etc/pam.d/acs-apg-common-banner|444
etc/pam.d/acs-apg-account-failure|444
etc/pam.d/acs-apg-account-success|444
etc/pam.d/acs-apg-auth-failure|444
etc/pam.d/acs-apg-auth-success|444
etc/pam.d/acs-apg-auth-role2group|444
etc/pam.d/acs-apg-lockout-tsadmin|444
etc/pam.d/acs-apg-lockout-tsgroup|444
etc/pam.d/acs-apg-password-local|444
etc/pam.d/acs-apg-session-heading|444
etc/pam.d/acs-apg-session-success|444
etc/smartd.conf|644
etc/ssh/sshd_config_4422|640
etc/ssh/sshd_config_mssd|640
etc/ssh/sshd_config_22|640
etc/ssh/sshd_config_830|640
etc/services_drbd|644
etc/services_md|644 
etc/systemd/system/dhcpd.service.d/lde.conf|644
etc/systemd/system/lde-sshd@sshd_config_22.service.d/apg_sshd.conf|644
etc/systemd/system/lde-sshd@sshd_config_4422.service.d/apg_sshd.conf|644
etc/systemd/system/lde-sshd@sshd_config_830.service.d/apg_sshd.conf|644
etc/systemd/system/lde-sshd@sshd_config_mssd.service.d/apg_sshd.conf|644
etc/sudoers.d/APG-comgroup_drbd|440
etc/sudoers.d/APG-comgroup_md|440
etc/sudoers.d/APG-tsadmin|440
etc/sudoers.d/APG-tsgroup_drbd|440
etc/sudoers.d/APG-tsgroup_md|440
etc/sysconfig/atftpd|644
etc/sysconfig/atftpd_AP2|644
etc/sysconfig/auditd|644
etc/sysconfig/openldap|644
etc/modprobe.d/lde-disable-unused-filesystems.conf|644
usr/lib/lde/config-management/apos_dhcpd-config|755
usr/lib/lde/config-management/apos_ntp-config|755
usr/lib/lde/config-management/apos_drbd-config|755
usr/lib/lde/config-management/apos_drbd-config_10g|755
usr/lib/lde/config-management/apos_drbd-config_VM|755
usr/lib/lde/config-management/apos_ftp-config|755
usr/lib/lde/config-management/apos_grub-config|755
usr/lib/lde/config-management/apos_ip-config|755
usr/lib/lde/config-management/apos_logrotd-config|755
usr/lib/lde/config-management/apos_rhosts-config|755
usr/lib/lde/config-management/apos_secacs-config|755
usr/lib/lde/config-management/apos_sshd-config|755
usr/lib/lde/config-management/apos_syslog-config|755
usr/lib/lde/config-management/apos_exports-config|755
usr/lib/lde/config-management/apos_rp-hosts-config|755
usr/lib/lde/failoverd-helpers/apg-defaults|755
usr/lib/systemd/scripts/apg-atftps.sh|755
usr/lib/systemd/scripts/apg-auditd.sh|755
usr/lib/systemd/scripts/apg-clearchipwdog.sh|755
usr/lib/systemd/scripts/apg-dhcpd.sh|755
usr/lib/systemd/scripts/apg-gid_checker.sh|755
usr/lib/systemd/scripts/apg-ldap.sh|755
usr/lib/systemd/scripts/apos-drbd.sh|755
usr/lib/systemd/scripts/apg-drbd-meta-convert|755
usr/lib/systemd/scripts/apos-finalize-system-conf.sh|755
usr/lib/systemd/scripts/apos-recovery-conf.sh|755
usr/lib/systemd/scripts/apos-system-conf.sh|755
usr/lib/systemd/system/apg-atftpd@.service|644
usr/lib/systemd/system/apg-atftps.service|644
usr/lib/systemd/system/apg-clearchipwdog.service|644
usr/lib/systemd/system/apg-dhcpd.service|644
usr/lib/systemd/system/apg-gid_checker.service|644
usr/lib/systemd/system/apg-ldap.service|644
usr/lib/systemd/system/apg-netconf-beep@.service|644
usr/lib/systemd/system/apg-netconf-beep.socket|644
usr/lib/systemd/system/apg-rsh@.service|644
usr/lib/systemd/system/apg-rsh.socket|644
usr/lib/systemd/system/apg-vsftpd@.service|644
usr/lib/systemd/system/apg-vsftpd.socket|644
usr/lib/systemd/system/apg-vsftpd-nbi@.service|644
usr/lib/systemd/system/apg-vsftpd-nbi.socket|644
usr/lib/systemd/system/apg-vsftpd-APIO_1@.service|644
usr/lib/systemd/system/apg-vsftpd-APIO_1.socket|644
usr/lib/systemd/system/apg-vsftpd-APIO_2@.service|644
usr/lib/systemd/system/apg-vsftpd-APIO_2.socket|644
usr/lib/systemd/system/apos-drbd.service|644
usr/lib/systemd/system/apos-early-system-config.service|644
usr/lib/systemd/system/apos-recovery-conf.service|644
usr/lib/systemd/system/apos-finalize-system-config.service|644
usr/lib/systemd/system/apos-system-config.service|644
usr/lib/systemd/system/auditd.service|644
usr/lib/systemd/system/dhcpd.service|644
usr/lib/systemd/system/smartd.service|644
usr/share/filem/internal_filem_root.conf|644
usr/share/filem/nbi_filem_root.conf|644
